# unidbg-fetch-qsign-bin/qsign@.service
[Unit]
Description=unidbg-fetch-qsign --basePath=txlib/%i
After=network.target

[Service]
WorkingDirectory=/opt/qsign
ExecStart=/opt/qsign/bin/unidbg-fetch-qsign --basePath=txlib/%i
Restart=always
RestartSec=1s

[Install]
WantedBy=multi-user.target
